# stt-bot.service
[Unit]
Description=STT Bot for Telegram (Voice Recognition)
After=network.target
Wants=network.target

[Service]
Type=simple
User=thar
WorkingDirectory=/home/thar/.openclaw/workspace
Environment=PATH=/home/thar/vosk_env/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/thar/.openclaw/workspace
ExecStart=/home/thar/vosk_env/bin/python /home/thar/.openclaw/workspace/stt_bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target